# etc/profile.d/antlr
#!/bin/bash

export LIB_PATH="$HOME/lib"
export CLASSPATH=".:$LIB_PATH"

alias antlr4='java -jar $LIB_PATH/antlr-4.7-complete.jar'
alias grun="java org.antlr.v4.gui.TestRig"
