# etc/profile.d/common
#!/bin/bash

export FZF_TMUX=1

alias vim="nvim"
