# etc/profile.d/docker
#!/bin/bash

export FUSION_MEMORY_SIZE=8192
export VIRTUALBOX_MEMORY_SIZE=8192
