# etc/profile.d/docker
#!/bin/bash

export FUSION_DISK_SIZE=40
export FUSION_MEMORY_SIZE=8192

export VIRTUALBOX_DISK_SIZE=40
export VIRTUALBOX_MEMORY_SIZE=8192
