# etc/profile.d/gcloud
#!/bin/sh

export GOOGLE_CREDENTIALS="$HOME/gcloud.json"
export GOOGLE_APPLICATION_CREDENTIALS="$HOME/gcloud.json"
