# etc/profile.d/go
#!/bin/bash

export GO111MODULE=on

if which go >&-
then
  export GOPATH="${HOME}/go"
  export PATH="${GOPATH}/bin:/usr/local/opt/go/libexec/bin:${PATH}"
fi

alias go-get="GO111MODULE=off go get"
