# etc/profile.d/hippo
#!/bin/bash

export PATH="$HOME/workspace/hippo/hippo-sdk/bin:$PATH"
