# etc/profile.d/http
#!/bin/bash

response_headers() {
  local url
  url="$1"
  shift

  curl -svL -o /dev/null "$url" "$@"
}
