# etc/profile.d/java
#!/bin/bash

export JAVA_HOME
export GRADLE_HOME

JAVA_HOME=$(/usr/libexec/java_home)
GRADLE_HOME="/usr/local/opt/gradle/libexec"
