# etc/profile.d/node
#!/bin/bash

export NVM_DIR="$HOME/.nvm"

# shellcheck disable=SC1091
source "/usr/local/opt/nvm/nvm.sh"
