# etc/profile.d/python
#!/bin/bash

export CFLAGS
export LDFLAGS

CFLAGS="-I$(brew --prefix openssl)/include"
LDFLAGS="-L$(brew --prefix openssl)/lib"

eval "$(pyenv init -)"
eval "$(pyenv virtualenv-init -)"
