# etc/profile.d/rlang
#!/bin/bash

export R_HOME="/usr/local/Cellar/r/3.4.1_2/lib/R"
