# etc/profile.d/rlang
#!/bin/bash

export R_HOME="/usr/local/opt/r/lib/R"
