# etc/profile.d/rust
#!/bin/bash

export RUST_SRC_PATH
export PATH=$HOME/.cargo/bin:$PATH

RUST_SYS_ROOT_PATH="$(rustc --print sysroot)"
RUST_SRC_PATH=$RUST_SYS_ROOT_PATH/lib/rustlib/src/rust/src
