# etc/profile.d/yarn
#!/bin/bash

if which yarn >&-
then
  export YARN_BIN
  YARN_BIN="$(yarn global bin)"
  export PATH="${PATH}:${YARN_BIN}"
fi
